# deploy/systemd/go-skeleton-worker.service
[Unit]
Description=go-skeleton Asynq worker
Documentation=https://github.com/myxiaoao/go-skeleton
After=network-online.target redis.service
Wants=network-online.target

[Service]
# notify: Worker sends sd_notify(READY=1) on startup + periodic WATCHDOG=1
# heartbeats. If a handler hangs (no panic, no error) systemd restarts the
# process per Restart= below; without watchdog those hangs would be silent
# until Asynqmon queue depth alerts catch up.
Type=notify
NotifyAccess=main
User=go-skeleton
Group=go-skeleton

WorkingDirectory=/opt/go-skeleton
EnvironmentFile=/etc/go-skeleton/.env
ExecStart=/opt/go-skeleton/bin/worker

Restart=on-failure
RestartSec=5s
# Worker does two-phase shutdown (Stop → Shutdown) on SIGTERM; give it
# room to finish in-flight tasks before systemd SIGKILLs it.
KillSignal=SIGTERM
TimeoutStopSec=45s
# Watchdog: same WATCHDOG_INTERVAL=10s heartbeat as the API. 60s is a
# wider margin than the API's 30s — worker tasks can legitimately block
# on long DB / external I/O between heartbeat ticks, and a false-positive
# restart mid-task is worse than for a stateless HTTP request.
WatchdogSec=60s
# Same rationale as API: HTTP/DB/Redis client pools + asynq client share
# FDs; default 1024 is tight under bursts.
LimitNOFILE=65535

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
