# deploy/systemd/go-skeleton-api.service
[Unit]
Description=go-skeleton HTTP API
Documentation=https://github.com/myxiaoao/go-skeleton
After=network-online.target postgresql.service redis.service
Wants=network-online.target

[Service]
# notify: 进程通过 sd_notify(READY=1) 告诉 systemd 启动完成；配合
# WatchdogSec= 让 systemd 在 watchdog 心跳缺失时按 Restart= 重启。
Type=notify
NotifyAccess=main
User=go-skeleton
Group=go-skeleton

# Layout assumed by docs/deploy.md:
#   /opt/go-skeleton/bin/{api,worker,migrate}
#   /etc/go-skeleton/.env   (config, mode 0640, owned by root:go-skeleton)
WorkingDirectory=/opt/go-skeleton
EnvironmentFile=/etc/go-skeleton/.env
ExecStart=/opt/go-skeleton/bin/api

Restart=on-failure
RestartSec=5s
# Send SIGTERM; the process handles graceful shutdown internally.
KillSignal=SIGTERM
TimeoutStopSec=30s
# Watchdog：30s 内没收到 WATCHDOG=1 就视为挂死并重启。
# 代码侧 WATCHDOG_INTERVAL=10s（≈ WatchdogSec/3），避免边界抖动。
WatchdogSec=30s
# 文件描述符上限：HTTP server + DB/Redis 连接池 + asynq client 累加；
# 默认 1024 在中等并发下吃紧。设 65535 给富余。
LimitNOFILE=65535

# Security hardening — sensible defaults for a stateless HTTP service.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

# Logs go to journald; use `journalctl -u go-skeleton-api -f` to follow.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
